# main.cron
# Запускаем поиск мошенников каждый день в 0 часов 13 минут
13 0 * * * home/demipt2/trof/main.py